# forex-monitor.service
[Unit]
Description=Forex Monitor — EUR/USD GBP/USD AUD/USD
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/forex-trading
ExecStart=/usr/bin/python3 /home/ubuntu/forex-trading/forex_monitor.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
